# docker/etc/services.d/shairport-sync/run
#!/usr/bin/execlineb -P
s6-setuidgid shairport-sync
/usr/local/bin/shairport-sync -vu --statistics